# TQ_IMX6Q_rootfs/lib/systemd/system/nvi.service
[Unit]
Description=Recover nvi edit sessions.

[Service]
Type=oneshot
ExecStart=/usr/share/vi/recover

[Install]
WantedBy=multi-user.target
